# netplay.service
[Unit]
Description=Sistema Netplay - Automação de Migração de Clientes
After=network.target
Wants=network.target

[Service]
Type=exec
User=netplay
Group=netplay
WorkingDirectory=/opt/netplay
Environment=PATH=/opt/netplay/venv/bin
EnvironmentFile=/opt/netplay/.env
ExecStart=/opt/netplay/venv/bin/gunicorn -c gunicorn.conf.py main:app
ExecReload=/bin/kill -HUP $MAINPID
KillMode=mixed
TimeoutStopSec=5
PrivateTmp=true
Restart=on-failure
RestartSec=5
StartLimitInterval=60s
StartLimitBurst=3

# Configurações de segurança
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=yes
ReadWritePaths=/opt/netplay
ReadWritePaths=/var/log/netplay
PrivateDevices=yes
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectControlGroups=yes

# Configurações de recursos
LimitNOFILE=65536
LimitNPROC=4096

# Logs
StandardOutput=append:/var/log/netplay/app.log
StandardError=append:/var/log/netplay/error.log
SyslogIdentifier=netplay

[Install]
WantedBy=multi-user.target